# scripts/com.suprawall.gtm.trigger.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.suprawall.gtm.trigger</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/alejandropeghin/Desktop/OpenClaw/LinkedIn-auto-research/scripts/autonomous_gtm_trigger.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/alejandropeghin/Desktop/OpenClaw/LinkedIn-auto-research</string>
    <key>StandardOutPath</key>
    <string>/Users/alejandropeghin/Desktop/OpenClaw/LinkedIn-auto-research/logs/cron.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/alejandropeghin/Desktop/OpenClaw/LinkedIn-auto-research/logs/cron.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
