# scripts/com.suprawall.gtm.dispatcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.suprawall.gtm.dispatcher.v2</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/alejandropeghin/Desktop/OpenClaw/LinkedIn-auto-research/scripts/launchd_dispatcher_wrapper.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>WorkingDirectory</key>
    <string>/Users/alejandropeghin/Desktop/OpenClaw/LinkedIn-auto-research</string>
    <key>StandardOutPath</key>
    <string>/Users/alejandropeghin/Desktop/OpenClaw/LinkedIn-auto-research/logs/cron.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/alejandropeghin/Desktop/OpenClaw/LinkedIn-auto-research/logs/cron.log</string>
</dict>
</plist>
